# activate-contact-site.service.example
[Unit]
Description=Activate Contact Site
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/activate-contact-site
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node /opt/activate-contact-site/server.js
Restart=always
RestartSec=3
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
